# backend/videogen.service
[Unit]
Description=VideoGen API Server
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/opt/videogen/backend
EnvironmentFile=/opt/videogen/backend/.env
ExecStart=/opt/videogen/backend/venv/bin/gunicorn server:app \
    -w 4 \
    -k uvicorn.workers.UvicornWorker \
    -b 0.0.0.0:8000 \
    --timeout 120 \
    --access-logfile /var/log/videogen/access.log \
    --error-logfile /var/log/videogen/error.log \
    --pid /run/videogen/server.pid
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
